# deploy/hincado-erp.service
[Unit]
Description=Hincado Directo ERP (Gunicorn)
After=network.target

[Service]
User=hincado
Group=hincado
WorkingDirectory=/opt/hincado-erp/interfaz_facturas
Environment="PATH=/opt/hincado-erp/interfaz_facturas/.venv/bin"
EnvironmentFile=/opt/hincado-erp/interfaz_facturas/.env
ExecStart=/opt/hincado-erp/interfaz_facturas/.venv/bin/gunicorn \
    --workers 2 \
    --bind 127.0.0.1:8000 \
    --timeout 120 \
    --access-logfile /var/log/hincado-erp/access.log \
    --error-logfile  /var/log/hincado-erp/error.log \
    backend:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
